# backend/unisub-backend.service
[Unit]
Description=UniSub Backend — 统一媒体订阅聚合器 API
After=network.target

[Service]
Type=simple
User=hermes
WorkingDirectory=/home/hermes/projects/unisub/backend
EnvironmentFile=/home/hermes/projects/unisub/.env
ExecStart=/home/hermes/.local/bin/uv run uvicorn app.main:app --host 0.0.0.0 --port 8002
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target